# packaging/systemd/user/battery-auditor-blackbox.service
[Unit]
Description=Battery Auditor black-box collector
Documentation=https://github.com/arfipod/battery-auditor
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/battery-auditor collect --mode blackbox --name blackbox-service
Restart=on-failure
RestartSec=5s
Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=6

[Install]
WantedBy=default.target
