# packaging/systemd/user/thinkpad-energy-manager-blackbox.service
[Unit]
Description=ThinkPad Energy Manager black-box collector
Documentation=https://github.com/arfipod/thinkpad-energy-manager
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/thinkpad-energy-manager collect --mode blackbox --name blackbox-service
Restart=on-failure
RestartSec=5s
Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=6

[Install]
WantedBy=default.target
