# packaging/systemd/user/thinkpad-energy-manager.service
[Unit]
Description=ThinkPad Energy Manager low-impact collector
Documentation=https://github.com/arfipod/thinkpad-energy-manager
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/thinkpad-energy-manager collect --mode diagnostic --name user-service
Restart=on-failure
RestartSec=5s
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
# Keep the collector quiet. It only reads sysfs and writes SQLite.

[Install]
WantedBy=default.target
